# custom.d/data/rc.local
#!/bin/sh

# Template start script
# This script will be copied to the SD card at /boot/rc.local
# On boot, it's called by /etc/rc.local


exit 0
