# copy_to_sd_boot/rc.local
#!/bin/sh

for f in /boot/rc.local.d/*.sh
do
  $f
done
